# crontab
@reboot chmod +x /var/www/ftpappend.sh
@reboot /var/www/ftpappend.sh
